# deploy/systemd/liquidity-migration-bybit-paper.service
[Unit]
Description=liquidity-migration paper (dry-run) event runner
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/liquidity-migration
EnvironmentFile=/etc/liquidity-migration/bybit-demo.env
Environment=DATA_ROOT=data/bybit-paper-event
Environment=CONFIG_PATH=configs/volume_alpha.default.yaml
Environment=SUBMIT_ORDERS=0
Environment=CONFIRM_DEMO_ORDERS=0
Environment=RECORD_DRY_RUN=1
Environment=TELEGRAM_ENABLED=0
Environment=STRATEGY_PROFILE=promoted
Environment=INTERVAL_SECONDS=60
Environment=USE_DAEMON=0
Environment=PYTHONDONTWRITEBYTECODE=1
Environment=MAX_ORDER_NOTIONAL_PCT_EQUITY=0
Environment=MAX_ENTRY_LAG_MINUTES=360
Environment=MAX_NEW_ENTRIES_PER_CYCLE=10
Environment=MAX_ACTIVE_SYMBOLS=3
Environment=UNIVERSE_RANK_END=400
Environment=UNIVERSE_MAX_SYMBOLS=400
Environment=UNIVERSE_MIN_TURNOVER_24H=0
Environment=WORKERS=8
Environment=ENTRY_LEVERAGE=2
Environment=WS_KLINES_ENABLED=1
ExecStart=/opt/liquidity-migration/scripts/run_bybit_demo_event_engine.sh
Restart=always
RestartSec=30
KillMode=control-group
# Bootstrap can take 60-120s on a cold start with a large kline universe.
# 90s default would SIGKILL the daemon if shutdown lands mid-bootstrap.
# Manager now respects shutdown_event, but keep margin for very slow REST.
TimeoutStopSec=180
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=false
ReadWritePaths=/opt/liquidity-migration

[Install]
WantedBy=multi-user.target
