# deploy/systemd/liquidity-migration-bybit-demo.service
[Unit]
Description=liquidity-migration Bybit demo event runner
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/liquidity-migration
EnvironmentFile=/etc/liquidity-migration/bybit-demo.env
Environment=DATA_ROOT=data/bybit-demo-event
Environment=CONFIG_PATH=configs/volume_alpha.default.yaml
Environment=SUBMIT_ORDERS=1
Environment=CONFIRM_DEMO_ORDERS=1
Environment=TELEGRAM_ENABLED=1
Environment=STRATEGY_PROFILE=promoted
Environment=INTERVAL_SECONDS=60
Environment=USE_DAEMON=1
Environment=PYTHONDONTWRITEBYTECODE=1
Environment=MAX_ORDER_NOTIONAL_PCT_EQUITY=0
# Was 360 (6h) when the feature pipeline could take 3-4h with the old
# REST-only kline path. WS klines deliver bars within ms of close and the
# daemon fires entries within ~20s of ready_ts (verified live 2026-05-25
# at Bybit createdTime = ready_ts + 20s). 180 (3h) is the right budget
# for the WS pipeline + occasional restart-bootstrap delays. Lowering
# further (60-120) risks dropping entries during deploy churn — see
# event_demo.py:71 for the alpha-decay rationale.
Environment=MAX_ENTRY_LAG_MINUTES=180
Environment=MAX_NEW_ENTRIES_PER_CYCLE=10
Environment=MAX_ACTIVE_SYMBOLS=3
# UNIVERSE_RANK_END / UNIVERSE_MAX_SYMBOLS = 0 → match-the-backtest mode:
# every active USDT-perp gets WS klines, daily-aggregated `liquidity_rank`
# is computed across the full Bybit perp universe (the same denominator
# the backtest uses), and the strategy's `universe_rank_max` filter
# applies on those daily ranks. Drops the live-ticker-turnover pre-filter
# that previously gave demo and backtest different rank denominators and
# therefore divergent entry selections on the same signal date (observed
# 2026-05-26 DRIFTUSDT: same data, demo entered, backtest rejected at
# rank-improvement-min=150 because prior7 was computed within 400 demo
# symbols vs 568 backtest symbols).
#
# To revert to the legacy narrow-universe demo (top-400 by ticker
# turnover, smaller kline store, but demo ≠ backtest), set both to 400
# and rebuild + restart.
Environment=UNIVERSE_RANK_END=0
Environment=UNIVERSE_MAX_SYMBOLS=0
Environment=UNIVERSE_MIN_TURNOVER_24H=0
Environment=WORKERS=2
Environment=ENTRY_LEVERAGE=2
Environment=WS_KLINES_ENABLED=1
ExecStart=/opt/liquidity-migration/scripts/run_bybit_demo_event_engine.sh
Restart=always
RestartSec=30
KillMode=control-group
# Bootstrap can take 60-120s on a cold start with a large kline universe.
# 90s default would SIGKILL the daemon if shutdown lands mid-bootstrap.
# Manager now respects shutdown_event, but keep margin for very slow REST.
TimeoutStopSec=180
Nice=5
IOSchedulingClass=best-effort
IOSchedulingPriority=7
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=false
ReadWritePaths=/opt/liquidity-migration

[Install]
WantedBy=multi-user.target
