# deploy/systemd/liquidity-migration-bybit-paper.service
[Unit]
Description=liquidity-migration paper (dry-run) event runner
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/liquidity-migration
EnvironmentFile=/etc/liquidity-migration/bybit-demo.env
Environment=DATA_ROOT=data/bybit-paper-event
Environment=CONFIG_PATH=configs/volume_alpha.default.yaml
Environment=SUBMIT_ORDERS=0
Environment=CONFIRM_DEMO_ORDERS=0
Environment=RECORD_DRY_RUN=1
Environment=TELEGRAM_ENABLED=0
Environment=STRATEGY_PROFILE=promoted
Environment=INTERVAL_SECONDS=60
# USE_DAEMON=1 since 2026-05-25 — was 0 (legacy bash per-cycle bootstrap)
# which added ~30-60s cycle latency. Daemon mode runs cycles in ~1-2s,
# matching demo's latency. Reconciliation pairing now sees both sides
# react to a signal within seconds, so MAX_ENTRY_LAG was lowered from
# 360 → 180 below. Paper doesn't submit (SUBMIT_ORDERS=0) so the WS
# trade client is unused; the WS klines + private state caches are
# shared design with the demo daemon.
Environment=USE_DAEMON=1
Environment=PYTHONDONTWRITEBYTECODE=1
Environment=MAX_ORDER_NOTIONAL_PCT_EQUITY=0
# Was 360 (6h) when paper ran legacy non-daemon mode with ~30-60s cycle
# latency. With paper now on USE_DAEMON=1 both daemons fire within seconds
# of ready_ts, so the lag budget can shrink. 180 (3h) keeps a comfortable
# buffer for occasional WS hiccups + restart bootstraps without chasing
# stale-alpha signals. Per event_demo.py docstring, 3h decay is acceptable.
Environment=MAX_ENTRY_LAG_MINUTES=180
Environment=MAX_NEW_ENTRIES_PER_CYCLE=10
# drop_all_4 promotion (2026-05-30): match the demo service — 12 slots (was 3),
# so the paper shadow tracks the validated package. See docs/preregistration/drop-all-4-promotion.md.
Environment=MAX_ACTIVE_SYMBOLS=12
# UNIVERSE_RANK_END / UNIVERSE_MAX_SYMBOLS = 0 → match-the-backtest mode.
# See liquidity-migration-bybit-demo.service for the rationale. The paper
# sleeve shadows the demo, so it MUST match identically; if demo is in
# match-the-backtest mode and paper isn't, demo↔paper slippage diagnostics
# would compare against an apples-to-oranges baseline.
Environment=UNIVERSE_RANK_END=0
Environment=UNIVERSE_MAX_SYMBOLS=0
Environment=UNIVERSE_MIN_TURNOVER_24H=0
Environment=WORKERS=2
Environment=ENTRY_LEVERAGE=2
Environment=WS_KLINES_ENABLED=1
ExecStart=/opt/liquidity-migration/scripts/run_bybit_demo_event_engine.sh
Restart=always
RestartSec=30
KillMode=control-group
# Bootstrap can take 60-120s on a cold start with a large kline universe.
# 90s default would SIGKILL the daemon if shutdown lands mid-bootstrap.
# Manager now respects shutdown_event, but keep margin for very slow REST.
TimeoutStopSec=180
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=false
ReadWritePaths=/opt/liquidity-migration

[Install]
WantedBy=multi-user.target
