# deploy/systemd/liquidity-migration-bybit-risk.service
[Unit]
Description=liquidity-migration fast Bybit demo exit risk runner
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/liquidity-migration
EnvironmentFile=/etc/liquidity-migration/bybit-demo.env
Environment=DATA_ROOT=data/bybit-demo-event
Environment=CONFIG_PATH=configs/volume_alpha.default.yaml
Environment=SUBMIT_ORDERS=1
Environment=CONFIRM_DEMO_ORDERS=1
Environment=TELEGRAM_ENABLED=1
Environment=ORDER_SUBMIT_MODE=ws_then_rest
Environment=REST_FALLBACK=1
Environment=PYTHONDONTWRITEBYTECODE=1
Environment=REST_RECONCILE_SECONDS=30
Environment=HEARTBEAT_SECONDS=10
Environment=FAST_EXECUTION_STREAM=0
# Combined-book deployment: set LONG_DATA_ROOT so this single risk service
# reads BOTH ledgers (short event_demo_* + long long_native_demo_*) and routes
# WS fills per the per-row sleeve column. Owner picked "extend ws_risk to
# handle both sides" over running two processes. Leave empty for short-only.
Environment=LONG_DATA_ROOT=data/bybit-long-demo-event
ExecStart=/opt/liquidity-migration/scripts/run_bybit_demo_ws_risk_engine.sh
Restart=always
RestartSec=2
KillMode=control-group
Nice=-5
IOSchedulingClass=best-effort
IOSchedulingPriority=4
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=false
ReadWritePaths=/opt/liquidity-migration

[Install]
WantedBy=multi-user.target
