# deploy/systemd/model050426-bybit-risk.service
[Unit]
Description=MODEL050426 fast Bybit demo exit risk runner
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/MODEL050426
EnvironmentFile=/etc/model050426/bybit-demo.env
Environment=DATA_ROOT=data/bybit-demo-event
Environment=CONFIG_PATH=configs/volume_alpha.default.yaml
Environment=SUBMIT_ORDERS=1
Environment=CONFIRM_DEMO_ORDERS=1
Environment=TELEGRAM_ENABLED=1
Environment=ORDER_SUBMIT_MODE=ws_then_rest
Environment=REST_FALLBACK=1
Environment=REST_RECONCILE_SECONDS=30
Environment=HEARTBEAT_SECONDS=10
Environment=FAST_EXECUTION_STREAM=0
ExecStart=/opt/MODEL050426/scripts/run_bybit_demo_ws_risk_engine.sh
Restart=always
RestartSec=2
KillMode=control-group
Nice=-5
IOSchedulingClass=best-effort
IOSchedulingPriority=4
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=false
ReadWritePaths=/opt/MODEL050426

[Install]
WantedBy=multi-user.target
