# deploy/systemd/model050426-bybit-canary.service
[Unit]
Description=MODEL050426 Bybit demo order-path canary
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=root
Group=root
WorkingDirectory=/opt/MODEL050426
EnvironmentFile=/etc/model050426/bybit-demo.env
Environment=DATA_ROOT=data/bybit-demo-event
Environment=CONFIG_PATH=configs/volume_alpha.default.yaml
Environment=SUBMIT_CANARY=1
Environment=CONFIRM_DEMO_ORDERS=1
Environment=CANARY_SYMBOL=DOGEUSDT
Environment=CANARY_SIDE=Buy
Environment=CANARY_NOTIONAL_USDT=8
Environment=CANARY_PRICE_DISTANCE_BPS=2000
ExecStart=/opt/MODEL050426/scripts/run_bybit_demo_canary.sh
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=false
ReadWritePaths=/opt/MODEL050426

[Install]
WantedBy=multi-user.target
