# deploy/systemd/model050426-bybit-demo.service
[Unit]
Description=MODEL050426 liquidity-migration Bybit demo event runner
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/MODEL050426
EnvironmentFile=/etc/model050426/bybit-demo.env
Environment=DATA_ROOT=data/bybit-demo-event
Environment=CONFIG_PATH=configs/volume_alpha.default.yaml
Environment=SUBMIT_ORDERS=1
Environment=CONFIRM_DEMO_ORDERS=1
Environment=TELEGRAM_ENABLED=1
Environment=STRATEGY_PROFILE=demo_relaxed
Environment=INTERVAL_SECONDS=60
Environment=MAX_ORDER_NOTIONAL_PCT_EQUITY=0
Environment=MAX_ENTRY_LAG_MINUTES=15
Environment=MAX_NEW_ENTRIES_PER_CYCLE=10
Environment=UNIVERSE_RANK_END=300
Environment=UNIVERSE_MAX_SYMBOLS=300
Environment=UNIVERSE_MIN_TURNOVER_24H=0
Environment=WORKERS=8
Environment=ENTRY_LEVERAGE=2
ExecStart=/opt/MODEL050426/scripts/run_bybit_demo_event_engine.sh
Restart=always
RestartSec=30
KillMode=control-group
Nice=5
IOSchedulingClass=best-effort
IOSchedulingPriority=7
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=false
ReadWritePaths=/opt/MODEL050426

[Install]
WantedBy=multi-user.target
