# deploy/systemd/model050426-bybit-canary.timer
[Unit]
Description=Run MODEL050426 Bybit demo order-path canary every 30 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=30min
RandomizedDelaySec=30s
Persistent=true
Unit=model050426-bybit-canary.service

[Install]
WantedBy=timers.target
